# modules/module-core/units/net/simpleiot/simpleiot.init
#!/sbin/openrc-run

description="Simple IoT — sensor data, telemetry, configuration, device management"

command="/usr/bin/siot"
command_background="yes"
pidfile="/run/simpleiot.pid"

depend() {
    need net
}
