# modules/module-core/units/net/docker-init/docker
#!/sbin/openrc-run

description="Docker daemon"

# Mirrors Alpine's docker-openrc package in spirit but only the bits yoe
# actually needs — most operators won't be tweaking dockerd.opts on a yoe
# image. /etc/conf.d/docker is honoured when present so this stays
# extensible without forking the script.

command="/usr/bin/dockerd"
command_args="${DOCKERD_OPTS}"
command_background="yes"
pidfile="/run/$RC_SVCNAME/pid"
start_stop_daemon_args="
    --pidfile $pidfile
    --background
    --make-pidfile
    --stdout /var/log/dockerd.log
    --stderr /var/log/dockerd.log
"

depend() {
    need cgroups net
    after firewall
    use dns
}

start_pre() {
    checkpath -d -m 0755 -o root:root /run/$RC_SVCNAME
    checkpath -f -m 0644 -o root:root /var/log/dockerd.log
}
