# modules/module-core/units/base/base-files/inittab.tmpl
# /etc/inittab — read by busybox init.
#
# Service supervision is delegated to OpenRC. busybox init handles PID 1
# duties (reaping, signals, getty respawn) and triggers the OpenRC
# runlevels in order: sysinit → boot → default.

::sysinit:/sbin/openrc sysinit
::sysinit:/sbin/openrc boot
::wait:/sbin/openrc default

{{.console}}::respawn:/sbin/getty -L {{.console}} 115200 vt100

::ctrlaltdel:/sbin/reboot
::shutdown:/sbin/openrc shutdown
